# deploy/murphys-bench-fetch-email.service
[Unit]
Description=Murphy's Bench — fetch inbound email (IMAP/POP3 -> tickets)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scs-tech
WorkingDirectory=/opt/murphys-bench
ExecStart=/opt/murphys-bench/venv/bin/python /opt/murphys-bench/manage.py fetch_inbound_email
