# deploy/murphys-bench-sla-check.service
[Unit]
Description=Murphy's Bench — check for overdue SLA tickets
After=network.target

[Service]
Type=oneshot
User=scs-tech
WorkingDirectory=/opt/murphys-bench
ExecStart=/opt/murphys-bench/venv/bin/python /opt/murphys-bench/manage.py check_sla_overdue
