# deploy/murphys-bench-backup.service
[Unit]
Description=Murphy's Bench nightly database backup
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=scs-tech
ExecStart=/opt/murphys-bench/scripts/backup_db.sh
