# deploy/murphys-bench-sla-check.timer
[Unit]
Description=Check Murphy's Bench SLA overdue tickets every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
